# config/service/sysSentry.service
[Unit]
Description=EulerOS System Inspection Frame
Requires=xalarmd.service sysSentry.socket
After=xalarmd.service sysSentry.socket

[Service]
ExecStart=/usr/bin/python3 /usr/bin/syssentry
ExecStop=/bin/kill $MAINPID
KillMode=mixed
Restart=on-failure
RestartSec=10s

[Install]
WantedBy=multi-user.target
