# src/libsentry/sh/log/logrotate-sysSentry.cron
#!/bin/sh

/usr/sbin/logrotate /etc/sysSentry/logrotate-sysSentry.conf -v -s /var/lib/logrotate-syssentry/logrotate.status
EXITVALUE=$?
if [ $EXITVALUE != 0 ]; then
	/bin/logger -t logrotate "ALERT exited abnormally with [$EXITVALUE], for details, see /var/log/sysSentry/logrotate.log"
fi
rm -rf /var/lib/logrotate-syssentry/logrotate.status
exit $EXITVALUE
